# system_files/usr/lib/systemd/system/systemd-suspend.service.d/10-armada-suspend.conf
[Unit]
Description=Armada Suspend

[Service]
# suspend-dispatch blocks until wake, so disable the start timeout.
ExecStart=
ExecStart=/usr/libexec/armada/suspend-dispatch
TimeoutStartSec=infinity
