# system_files/usr/lib/systemd/system/armada-installer-visibility.service
[Unit]
Description=armada installer KDE entry visibility
After=systemd-tmpfiles-setup.service
Before=display-manager.service
RequiresMountsFor=/run

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/armada/installer-visibility
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
