# system_files/usr/lib/systemd/system/armada-session-default.service
[Unit]
Description=Reset session autologin to gamemode each boot
Before=display-manager.service
After=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/armada/session-control default-gamemode
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
