# system_files/usr/lib/systemd/system/armada-power-switch.service
[Unit]
Description=armada power-button + lid watcher (fake suspend)
After=multi-user.target

[Service]
ExecStart=/usr/libexec/armada/power-switch-watch
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
